# docker-assets/webapp/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
autorestart=true
directory=/var/app/current
command=node .
